# go-cronjob/cron.d/root
* * * * * /app/script/job_script.sh >> /app/log/job_script.log 2>&1
